# bluetooth/bt-auto-connect.service
[Unit]
Description=Auto-connect Bluetooth headset at login
After=default.target bluetooth.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/bt-auto-connect.sh

[Install]
WantedBy=default.target
